# admin/roles/frontend/files/supervisor.conf
[program:frontend]
command=/home/james/frontend/bin/frontend
directory=/home/james/frontend
priority=200
stopsignal=INT
user=james

